# mkdocs.conf
[program:mkdocs]
command=mkdocs serve -f /var/www/mkdocs.yml -a 0.0.0.0:80
autostart=false
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
